# northwoodlogic-overlay/etc/inittab
::sysinit:/sbin/startup
::shutdown:/bin/umount -a -r
::askfirst:-/bin/cttyhack login -f root
::respawn:/bin/runsvdir -P /var/service
